# debian/gatekeeper.cron.daily
#!/bin/sh

/usr/bin/find /var/log/gatekeeper \
  -maxdepth 1                     \
  -type f                         \
  -name gatekeeper_\*.log         \
  -mtime +1                       \
  -delete
